# setup/video_system/rpi/config/cfgtopcopy/etc/systemd/system/systemd-remount-fs.service.d/override.conf
# Override to prevent remount errors when using overlay filesystem
# Overlay filesystems don't support remounting, so we skip this service
[Unit]
ConditionPathExists=!/sys/kernel/debug/overlayfs
