# setup/video_system/rpi/config/cfgtopcopy/etc/systemd/system/vlc.service
[Unit]
Description=VLC Media Player (Kiosk Mode)
Documentation=man:vlc(1)
After=openbox.service
Requires=openbox.service
PartOf=graphical.target

[Service]
Type=simple
User=ccpadmin
Group=ccpadmin
# Load X11 environment
EnvironmentFile=/config/xorg.env
Environment=DISPLAY=:0
# Disable D-Bus and PulseAudio (we use ALSA directly)
Environment=DBUS_SESSION_BUS_ADDRESS=
Environment=PULSE_SERVER=
# Create XDG_RUNTIME_DIR if systemd didn't (fallback for user services)
ExecStartPre=/bin/sh -c 'mkdir -p /run/user/$(id -u ccpadmin) && chown ccpadmin:ccpadmin /run/user/$(id -u ccpadmin)'
ExecStart=/config/start_vlc_v2.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=false

[Install]
WantedBy=graphical.target
